# crontab
#!/bin/bash

find /home/photo -iname "*.jpg" -exec jhead -n%Y%m%d-%H%M%S {} \+
